# deploy/irc-core.service
# irc-core systemd service — the persistent IRC connection daemon
# Place at /etc/systemd/system/irc-core.service
#
# irc-core owns the raw IRC TCP/TLS sockets so that redeploying/restarting
# cryptirc.service (the frequent case) no longer drops any user's IRC
# connection. It's a thin "socket-keeper" — cryptirc.service still does 100%
# of parsing/logging/E2E/push; this daemon only dials, registers, and holds
# the wire open, handing lines back and forth over a local Unix socket.
#
# Install:
#   sudo cp deploy/irc-core.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now irc-core
#
# Start this BEFORE cryptirc.service on a fresh boot (see After= below) — a
# cold-started daemon just waits idle for Dial requests, so starting it first
# is always safe. Ordinary code-only redeploys of cryptirc.service should NOT
# restart this unit — see deploy/update.sh's --restart-daemon flag.

[Unit]
Description=irc-core — persistent IRC connection daemon for CryptIRC
Documentation=https://github.com/gh0st68/CryptIRC
After=network-online.target
Wants=network-online.target
# Ordering only (no Requires=) — cryptirc.service should start after this,
# but this daemon crashing/restarting must never cascade-stop the web
# service, and vice versa. Lifecycle-decoupling is the entire point.
Before=cryptirc.service
# NEVER latch into a permanent `failed` state. The old StartLimitIntervalSec=60 +
# StartLimitBurst=5 meant 5 fast crashes (RestartSec=5 → ~25s, inside the 60s window)
# would burn all tokens and leave the daemon dead until a human ran `reset-failed`.
# On a never-touch box that is the exact catastrophe this deployment exists to avoid.
# `=0` disables the rate limiter → systemd retries forever; a manual `systemctl stop`
# still stops it cleanly, so this is safe.
StartLimitIntervalSec=0

[Service]
# Type=notify so the daemon can drive the systemd WATCHDOG (see WatchdogSec below) —
# the ONLY backstop for a hung-but-alive daemon (write-stall, dial-stall, a deadlock),
# and it cannot be retrofitted after the binary is frozen.
Type=notify
NotifyAccess=main
# If the daemon stops sending WATCHDOG=1 keepalives for this long, systemd kills and
# restarts it. The daemon heartbeats only while the async RUNTIME is demonstrably
# responsive (its watchdog task is scheduled AND the central lock is acquirable), so a
# wedged executor / all-threads-blocked / deadlock is caught here. Per-loop logic hangs
# are designed out upstream with hard timeouts. Heartbeat cadence is ~WatchdogSec/3.
WatchdogSec=60
User=cryptirc
Group=cryptirc
WorkingDirectory=/opt/cryptirc

# The compiled daemon binary (built as `irc_core` — see src/bin/irc_core.rs)
ExecStart=/opt/cryptirc/irc_core

# Graceful shutdown — sends SIGTERM, waits up to 30s, then SIGKILL
ExecStop=/bin/kill -s TERM $MAINPID
TimeoutStopSec=30
KillMode=mixed

# Restart policy — always restart (a patiently-retrying daemon beats a permanently-dead
# one). `Restart=always` does NOT restart on a clean `systemctl stop`, so intentional
# stops still work. Paired with StartLimitIntervalSec=0 above → truly never gives up.
Restart=always
RestartSec=5

# Resource ceilings so a runaway (leaked FDs, a memory-accreting bug over years of
# uptime) fails FAST into a clean restart instead of the OS-OOM-killer lottery that
# could take neighbouring services with it.
MemoryMax=512M
TasksMax=256

# Files/dirs the daemon creates (the IPC socket) default to 0600/0700.
UMask=0077

# Same data dir as cryptirc.service — the daemon derives its Unix socket
# path as $CRYPTIRC_DATA/irc-core.sock (override with CRYPTIRC_IPC_SOCK if
# ever needed). It never reads/writes any other file in this directory:
# NetworkConfig secrets are decrypted by the web process, never the daemon.
Environment=CRYPTIRC_DATA=/var/lib/cryptirc
Environment=RUST_LOG=info

# ── Security hardening (matches cryptirc.service's posture) ────────────────

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/lib/cryptirc
ProtectHome=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
# AF_INET/6 for the outbound IRC TCP/TLS dial, AF_UNIX for the IPC socket.
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
CapabilityBoundingSet=
AmbientCapabilities=
LockPersonality=true
ProtectClock=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectControlGroups=true
ProtectHostname=true

# Generous FD headroom. Real usage is tiny (a few dozen IRC sockets + one IPC client),
# but a transient FD spike must never trip EMFILE and cascade through accept(). Well
# above any realistic need so the accept loop's error-swallowing is belt-and-suspenders.
LimitNOFILE=16384
LimitNPROC=64

[Install]
WantedBy=multi-user.target
